# forge/acquisition/systemd/units/forge-acquisition-initialize.target
[Unit]
Description=Forge acquisition initialized
Before=forge-acquisition-control.target forge-acquisition-start.target
After=forge-acquisition-stop.target
Conflicts=forge-acquisition-stop.target